# deploy/alice/wonderlands-app.service
[Unit]
Description=Alice Wonderlands app API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=alice
Group=alice
WorkingDirectory=/home/alice/wonderlands/current/apps/server
EnvironmentFile=/home/alice/wonderlands/shared/server.env
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run start
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
